# backend/docker.cron
* * * * * python /var/www/metcalf/backend/manage.py refresh_datafeeds 2>&1 | /var/log/cron.log -t django-cron
@daily python /var/www/metcalf/backend/manage.py schedule_datafeeds -a 2>&1 | /var/log/cron.log -t django-cron