# docker-files/crontab
*/1 * * * * python3 manage.py refresh_datafeeds
@daily python3 manage.py schedule_datafeeds --all 